# supervisor.conf
[program:echo]
command = /tmp/echo.sh
#command = /Users/dabou/MyProjects/odo-supervisor/echo.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_events_enabled=true