# supervisor.conf
[program:echo]
command = /tmp/echo.sh
stdout_logfile=/dev/stdout
stdout_events_enabled=true