# supervisor-local.conf
[program:echo]
command = ./echo.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_events_enabled=true