# conf/supervisor.conf
[program:echo]
command = /tmp/echo.sh
stdout_logfile=/dev/stdout
stdout_events_enabled=true

[inet_http_server]
port=localhost:9001